# data/danctnix-tweakd.initd.in
#!/sbin/openrc-run

name=$RC_SVCNAME
description="DanctNIX tweaks daeamon"
command="/usr/bin/danctnix-tweakd"
