# data/postmarketos-tweakd.initd.in
#!/sbin/openrc-run

name=$RC_SVCNAME
description="postmarketOS tweaks daeamon"
command="/usr/bin/pmos-tweakd"
